# flood.conf
[program:flood]
command=/usr/bin/python3 /home/ftc/flood/index.py
directory=/home/ftc/flood
user=ftc
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

